# system_stuff/add_ips.service
[Unit]
Description=Add two VIPs for containers
After=network.target

[Service]
Type=onshot
ExecStart="ifconfig wlan0:1 192.168.50.242"
ExecStart="ifconfig wlan0:2 192.168.50.243"

[Install]
WantedBy=multi-user.target

