# system_stuff/add_ips.service
[Unit]
Description=add two IPs
After=remote-fs.target

[Service]
Type=oneshot
ExecStart=ifconfig wlan0:1 192.168.42.242
ExecStart=ifconfig wlan0:2 192.168.42.243
ExecStart=ifconfig wlan0:3 192.168.42.2


[Install]
WantedBy=multi-user.target

