# system_stuff/restart_dragon.service
[Unit]
Description=restart dragon containers
After=multi-user.target

[Service]
Type=oneshot
ExecStart=docker restart codeadmin
ExecStart=docker restart dragon-ssh
ExecStart=docker restart dragon-apache
ExecStart=docker restart router-apache

[Install]
WantedBy=multi-user.target

